# module/Analytics/config/crontab
@weekly www-data php /vagrant/public/index.php analytics_weekly_aggregation 2>&1