# module/Analytics/config/crontab
@weekly php /vagrant/public/index.php analytics_weekly_aggregation